# bin/crontab
* * * * * /usr/sbin/screencapture /Users/yisheng/Dropbox/grepawk/screenshot.jpg  2>&1

